# deploy/omniroute.service
[Unit]
Description=OmniRoute — KiroAI Claude API proxy
After=network.target

[Service]
Type=simple
User=root
Environment=HOME=/root
WorkingDirectory=/root
ExecStart=/usr/bin/node /usr/bin/omniroute
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omniroute

[Install]
WantedBy=multi-user.target
